# crontab
00 21 * * * root curl -X GET http://localhost:3000/api/check-schedules
# Don't remove the empty line at the end of this file. It is required to run the cron job
